# deployment/systemd/portfolio-v2-backend.service
[Unit]
Description=Portfolio V2 Backend API
Documentation=https://go.dev/doc/
After=network.target redis.service

[Service]
Type=simple
User=serkanursavas
WorkingDirectory=/home/serkanursavas/projects/portfolio-v2/backend
Environment=GIN_MODE=release
Environment=REDIS_URL=redis://localhost:6379
Environment=PORT=8082
EnvironmentFile=/home/serkanursavas/projects/portfolio-v2/backend/.env
ExecStart=/home/serkanursavas/projects/portfolio-v2/backend/portfolio-backend
Restart=always
RestartSec=10

# Build command (run manually): go build -o portfolio-backend main.go

# Logging - using journal instead of obsolete syslog
StandardOutput=journal
StandardError=journal

# Security - Relaxed for debugging
NoNewPrivileges=no
# ProtectHome=yes
# ProtectSystem=strict
# ReadWritePaths=/home/serkanursavas/projects/portfolio-v2/backend/uploads
# ReadWritePaths=/home/serkanursavas/projects/portfolio-v2/backend/blog-upload
# ReadWritePaths=/home/serkanursavas/projects/portfolio-v2/backend/skills-upload

[Install]
WantedBy=multi-user.target